# deployment/systemd/alpha-score-backend.service
[Unit]
Description=Alpha-Score Backend API Service
After=network.target

[Service]
Type=simple
User=xd
Group=xd
WorkingDirectory=/home/xd/project/alpha-score/backend

# 环境变量
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/xd/project/alpha-score/backend/venv/bin:/usr/local/bin:/usr/bin"

# 启动命令
ExecStart=/home/xd/project/alpha-score/backend/venv/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65535
MemoryLimit=1G
CPUQuota=200%

# 安全加固
NoNewPrivileges=true
PrivateTmp=true

# 标准输出和错误输出
StandardOutput=append:/home/xd/project/alpha-score/logs/systemd/backend-stdout.log
StandardError=append:/home/xd/project/alpha-score/logs/systemd/backend-stderr.log

[Install]
WantedBy=multi-user.target
